# automata-app/tools/automata.service
[Unit]
Description=Orcun's home automation service.

[Service]
Type=simple
ExecStart=/usr/local/bin/automata

[Install]
WantedBy=multi-user.target
